# Qwen_llama/motia/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:iii-engine]
command=/root/.local/bin/iii -c /app/motia/iii-config.yaml
directory=/app/motia
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:iii-console]
command=/root/.local/bin/iii-console --enable-flow --host 0.0.0.0 --port 3113 --engine-host localhost --engine-port 3121 --ws-port 49134 --bridge-port 49134
directory=/app/motia
startdelay=3
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0